# claude-bot.service
[Unit]
Description=Claude Telegram Bot
After=network.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/opt/claude-telegram-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/claude-telegram-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
